# runit/calendar-api/run
#!/bin/sh
set -eu

set -a
. /etc/calendar-api/env
set +a

exec 2>&1
exec /usr/local/bin/calendar-api
